# deployment/domain-monitor.service
[Unit]
Description=Domain Monitor Service
After=network.target

[Service]
Type=simple
User=monitor
Group=monitor
WorkingDirectory=/opt/domain-monitor
Environment="PATH=/opt/domain-monitor/venv/bin"
ExecStart=/opt/domain-monitor/venv/bin/python /opt/domain-monitor/src/main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/domain-monitor/output.log
StandardError=append:/var/log/domain-monitor/error.log

[Install]
WantedBy=multi-user.target